# etc/systemd/gofer-https.socket
[Socket]
ListenStream=443
FileDescriptorName=https
Service=gofer.service

