# etc/systemd/gofer-http.socket
[Socket]
ListenStream=80
FileDescriptorName=http
Service=gofer.service

